# stemcell_builder/stages/system_oracle/assets/etc/update-motd.d/51-cloudguest
CLOUD_IMG: This file was created/modified by the Cloud Image build process
# This file is not managed by a package.  If you no longer want to
# see this message you can safely remove the file.
echo ""
echo "  Get cloud support with Ubuntu Advantage Cloud Guest:"
echo "    http://www.ubuntu.com/business/services/cloud"
