# stemcell_builder/stages/base_opensuse/assets/runit.service
[Unit]
Description=Runit Process Supervisor

[Service]
ExecStart=/sbin/runsvdir-start
Restart=always
TasksMax=infinity

[Install]
WantedBy=multi-user.target
